# systemd/crypto-recorder.service
[Unit]
Description=Binance Market Data Recorder (Deterministic Native)
Documentation=file:///home/zsom/services/CryptoRecorder/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zsom
WorkingDirectory=/home/zsom/services/CryptoRecorder
Environment="CRYPTO_RECORDER_TOP_SYMBOLS=50"
ExecStart=/home/zsom/services/CryptoRecorder/.venv/bin/python recorder.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto-recorder

# Resource limits
MemoryMax=8G
CPUQuota=50%

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
